# python/openclaw.ini
[program:openclaw]
process_name=%(program_name)s_%(process_num)02d
command=openclaw gateway --port 18789 --verbose
autostart=true
autorestart=true
priority=1
pidfile=/tmp/openclaw.pid
redirect_stderr=true
stdout_logfile=/var/www/logs/process-openclaw.log
environment=HOME="/root",NODE_ENV="production"
