# python/socat-relay.ini
[program:socat-codex-oauth]
command=socat TCP-LISTEN:41455,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:1455
autostart=true
autorestart=true
priority=0
redirect_stderr=true
stdout_logfile=/var/www/logs/socat-codex-oauth.log